# cmd/gnmi_collector/docker/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log                ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; do not daemonize
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/                           ; where child log files will live

[program:collector]
command = /gnmi_collector -config_file /config/example.cfg -cert_file /config/cert.pem -key_file /config/key.pem -port 8888 -alsologtostderr
startsecs=10
startretries=10
autorestart=true
